# systemd/damnit-services.target
[Unit]
Description=DAMNIT services

[Install]
WantedBy=default.target
